# nifr/playbooks/devcontainer/templates/s6/redis-monitor/services.d/redis-monitor/run
#!/usr/bin/env sh

exec s6-setuidgid redis \
  redis-cli --csv monitor \
  2>&1
