# nifr/playbooks/devcontainer/templates/s6/redis-monitor/services.d/redis-monitor/run
 #!/usr/bin/with-contenv sh

exec 2>&1
exec s6-setuidgid 'redis' \
  sh -c 'redis-cli --csv -h "${REDIS_HOST:-127.0.0.1}" -p "${REDIS_PORT:-6379}" -n "${REDIS_DB:-0}" monitor 2>&1 | grep -ve ping'
