# nifr/playbooks/devcontainer/templates/s6/mysqld/services.d/mysqld/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid mysql \
  mysqld \
    --defaults-file=/etc/mysql/my.cnf \
  2>&1
