# nifr/playbooks/devcontainer/templates/s6/minio/services.d/minio/run
#!/usr/bin/with-contenv sh

# see:
exec s6-setuidgid minio \
  minio server /opt/minio/data \
    --console-address "127.0.0.1:9001" \
    --json \
    2>&1
