# nifr/playbooks/devcontainer/templates/s6/minio/services.d/minio/run
#!/usr/bin/with-contenv sh

exec 2>&1
exec \
  MINIO_ROOT_USER=admin \
  MINIO_ROOT_PASSWORD=password \
  minio server /data \
    --console-address ":9001" \
    --json
