# nifr/playbooks/devcontainer/templates/s6/mysql-import/services.d/mysql-import/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid mysql \
  echo 'I am just running once!'
