# nifr/playbooks/devcontainer/templates/s6/grafana-agent/services.d/grafana-agent/run
#!/usr/bin/env sh

cd '{{ grafana_agent.home }}'

# see: https://grafana.com/docs/agent/latest/configuration/
#
# ## remote config (beta)
# -enable-features=remote-configs \
# -config.file=https://example.com/config.yml \
# -config.url.basic-auth-user=<user> \
# -config.url.basic-auth-password-file=<file>
#
exec s6-setuidgid '{{ grafana_agent.user }}' \
  grafana-agent \
    -config.file '{{ grafana_agent.config_file }}' \
    -config.expand-env \
    2>&1
