# nifr/playbooks/devcontainer/templates/s6/grafana-agent/services.d/grafana-agent/run
#!/usr/bin/env sh

cd '{{ grafana_agent.home }}'

exec s6-setuidgid '{{ grafana_agent.user }}' \
  grafana-agent \
    -config.file '{{ grafana_agent.config_file }}' \
    -config.expand-env \
    2>&1
