# nifr/playbooks/devcontainer/templates/s6/openresty/services.d/openresty/run
#!/usr/bin/execlineb -P

# s6-setuidgid openresty
fdmove -c 2 1
openresty
